# night_shutdown/pishutdown.service
[Service]
ExecStart=/usr/bin/python /home/pi/night_shutdown/pishutdown.py
WorkingDirectory=/home/pi/night_shutdown/
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pishutdown
User=root
Group=root

[Install]
WantedBy=multi-user.target
